# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/usr/sbin/openvpn --client --dev tun0 --config /etc/openvpn/client.conf

[program:syslog]
command=/usr/sbin/rsyslogd -n
